# docker/production/backup/config/crontab/localserver
SHELL=/bin/sh

0 * * * *		rclone sync ~/backups/ cc:"/ModelIt/Database Backups"
